# config/crontab
@reboot /src/nest_push.py && date >> /var/log/cron.log
*/1 * * * * /src/nest_push.py && date >> /var/log/cron.log
